# claude-session.service
[Unit]
Description=Claude Code persistent tmux session
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
User=harsh
WorkingDirectory=%h
Environment=HOME=%h
Environment=USER=harsh
Environment=NODE_ENV=production

# Kill any stale tmux session with same name before starting
ExecStartPre=-/usr/bin/tmux kill-session -t claude-vps 2>/dev/null
ExecStartPre=/bin/sleep 1

# Start new detached tmux session running Claude
# --dangerously-skip-permissions is required for autonomous operation
ExecStart=/usr/bin/tmux new-session -d -s claude-vps \
    /bin/bash -c 'while true; do \
        claude --dangerously-skip-permissions 2>&1 | tee -a ~/.claude-persist.log; \
        echo "[$(date)] Claude exited (code $?) — restarting in 5s..." >> ~/.claude-persist.log; \
        sleep 5; \
    done'

# Never give up restarting
Restart=always
RestartSec=10

# Auto-restart even if killed manually
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
